# security/ulogd/init.d/ulogd
#!/bin/bash
#
# chkconfig: 345 81 19
# description: ulogd is the userspace logging daemon for netfilter/iptables
# processname: ulogd
#

ARGS="--daemon"

start () {
  get_data
  echo -n "Starting $SHORT: "
  ulogd $ARGS &> /dev/null
  if [ $? == 0 ] ; then
    mkpid
    echo -e $RESULT_OK
  else
    echo -e $RESULT_FAIL
  fi
}  

. /lib/lsb/init-functions

