# crypto/heimdal/systemd.d/heimdal-kdc.service
[Unit]
Description=Heimdal Kerberos 5 Server
After=syslog.target network.target

[Service]
ExecStart=/usr/heimdal/libexec/kdc
Restart=always

[Install]
WantedBy=multi-user.target
