# audio/audacity/profile.d/audacity.rc
if [ ! -f ~/.Audacity ]; then
exec 5>~/.Audacity
cat >&5 <<EOF
Version=1.1.1
PrefsVersion=1.1.1r1
EOF
exec 5>&-
fi
