# net/ypserv/init.d/ypserv
#!/bin/sh
#
# ypserv:      Starts the yp-server
#
# Author:      kongar, 20020417
#
# chkconfig: 2345 16 84
#
# description: this script contains a trap to catch the unsuccessful ypserv \
#              execution. Since ypserv returns NOTHING if not successful, \
#              this seems to be a necessity. \
#              This script seems to be tolerate to its name
#
# processname: ypserv
# pidfile: /var/run/ypserv.pid


YPSERV_PID=/var/run/ypserv.pid

case  $1  in
  start|restart)  echo "$1ing YP server"
                  [ -e "$YPSERV_PID" ] && \
                    kill `cat $YPSERV_PID` 2>/dev/null && \
                    sleep 1
                  /usr/sbin/ypserv && sleep 1
                  RES=`pidof /usr/sbin/ypserv`
                  if [ -z "$RES" ]; then
                     /usr/sbin/ypserv -d&
                     rm $YPSERV_PID
                  fi
                  ;;

           stop)  echo "$1ping YP server"
                  [ -e "$YPSERV_PID" ] && kill `cat $YPSERV_PID` 2>/dev/null
                  ;;

              *)  echo "Usage: $0 {start|stop|restart}"
                  ;;
esac
