# graphics/sane-backends/systemd.d/saned.socket
[Unit]
Description=saned incoming socket

[Socket]
Accept=yes
ListenStream=6566
SocketMode=0775

[Install]
WantedBy=sockets.target
