# devel/binutils-sparc64/profile.d/binutils-sparc64.rc
export PATH=/opt/lunar/binutils-sparc64/bin:$PATH
