# sql/mongodb/init.d/mongodb
#!/bin/bash
#
# Startup script for mongodb
#
# chkconfig: 345 90 20
# description: A scalable, high-performance, key-value store, document-oriented database
# processname: mongod
# pidfile: /var/run/mongodb.pid

# config settings
[ -f /etc/config.d/mongodb ] && . /etc/config.d/mongodb

# MongoDB config
[ -f /etc/mongodb.conf ] || exit 1
[ -x /usr/bin/mongod ] || exit 1

ARGS=""

MDBPID=$(cat /var/run/mongodb.pid)

if [ -z "$MDB_USER" ] ; then
	echo "ERROR: A MDB_USER must be set!"
	exit 1
fi

start()
{
  echo -n "Starting mongodb: "
  if [ -n "$MDBPID" ]; then
    echo "already running."
  else
    su - $MDB_USER -c "/usr/bin/mongod --config /etc/mongodb.conf --pidfilepath /var/run/mongodb.pid &> /dev/null" &&
    echo -e $RESULT_OK || echo -e $RESULT_FAIL
  fi
}

stop()
{
  echo -n "Stopping mongodb: "
  if [ -z "$MDBPID" ]; then
    echo "not running."
  else
    kill $MDBPID &&
    echo -e $RESULT_OK || echo -e $RESULT_FAIL
  fi
}

restart()
{
  $0 stop
  sleep 1
  $0 start
}
. /lib/lsb/init-functions
