# crypto/heimdal/xinetd.d/krb5-telnet
# description: The kerberized telnet server accepts normal telnet sessions, \
#              but can also use Kerberos 5 authentication.
service telnet
{
        flags           = REUSE
        socket_type     = stream        
        wait            = no
        user            = root
        server          = /usr/libexec/telnetd
        server_args     = -a user
        log_on_failure  += USERID
        disable         = no
}
