# crypto/heimdal/xinetd.d/kx
service kx
{
        type            = UNLISTED
        flags           = REUSE
        socket_type     = stream        
        wait            = no
        user            = root
        server          = /pkg/heimdal/default/libexec/kxd
        log_on_failure  += USERID
        disable         = no
        port            = 2111
}
