# eclipse/eclipse-cpp/profile.d/eclipse-cpp.rc
export PATH=${PATH}:/opt/lunar/eclipse-cpp
