# compilers/go/profile.d/go.rc
export GOROOT="/usr/lib/go"
export GOPROXY=https://proxy.golang.org
