# net/ypbind-mt/init.d/ypbind
#!/bin/sh
#
# description: ypbind
# chkconfig: 345 80 30
# processname: ypbind

PID_FILE=/var/run/ypbind.pid
NISDOMAIN=`cat  /etc/nisdomain`

case $1 in
  start|restart)  echo   "$1ing ypbind."
                  ps  -C  portmap  >  /dev/null  ||  /etc/init.d/portmap start
                  pkill -f "\/usr\/sbin\/ypbind$"     &&  sleep 1
                  pkill -9 -f "\/usr\/sbin\/ypbind$"     &&  sleep 1
                  /bin/nisdomainname $NISDOMAIN
                  /usr/sbin/ypbind
                  ;;
              
           stop)  echo   "$1ping ypbind."

                  pkill -f "\/usr\/sbin\/ypbind$"     &&  sleep 1
                  pkill -9 -f "\/usr\/sbin\/ypbind$" 
                  [ -e "$PID_FILE" ] && rm $PID_FILE
                  ;;

              *)  echo  "Usage: $0 {start|stop|restart}"
                  ;;
esac
