# zdeprecated/egcs64/profile.d/egcs64.rc
export PATH=/opt/lunar/egcs64/bin:$PATH
