# zbeta/x11/x11-proto/xproto/profile.d/x11r7.rc
export PKG_CONFIG_PATH=/usr/X11R7/lib/pkgconfig:$PKG_CONFIG_PATH
