# x11-apps/openoffice-bin/profile.d/openoffice-bin.rc
# openoffice-bin specific stuff

[ -d /opt/lunar/openoffice/program ] &&
    PATH="$PATH:/opt/lunar/openoffice/program"

