# x11-apps/openoffice-bin/profile.d/openoffice-bin.rc
# openoffice-bin specific stuff

[ -d /opt/lunar/openoffice/openoffice.org3/program ] &&
    PATH="$PATH:/opt/lunar/openoffice/openoffice.org3/program"

