# java/clojure/profile.d/clojure.rc
export CLOJURE_HOME=/usr/share/java/clojure
export CLASSPATH=$CLASSPATH:$CLOJURE_HOME/clojure.jar
