# science/alliance/profile.d/alliance.rc
export PATH=$PATH:/opt/lunar/alliance/bin
export MANPATH=$MANPATH:/opt/lunar/alliance/man
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/opt/lunar/alliance/lib
. /opt/lunar/alliance/etc/alc_env.sh
