# net/finger/xinetd.d/finger
# default: on
# description: The finger server answers finger requests. Finger is \
#	a protocol that allows remote users to see information such \
#	as login name and last login time for local users.
service finger
{
	socket_type	= stream
	wait		= no
	user		= nobody
	server		= /usr/sbin/in.fingerd
	disable		= yes
}
