# audio/alsa-utils/init.d/alsasettings
#!/bin/bash
#
# alsasettings  Takes care of Storing the ALSA mixer values
#		2.6.x kernel systems
#
# chkconfig: 345 99 12
# description: Store ALSA settings on 2.6.x based systems

PATH="$PATH:/usr/sbin:/usr/bin"

function store() {
  echo "Storing ALSA mixer settings..."
  alsactl store
  amixer set Master mute &> /dev/null
}

# Added for restoring alsa mixer levels. neocrust@gmail.com
function restore() {
 echo "Restoring ALSA mixer settings..."
 alsactl restore
}

case "$1" in
  start)
	restore
	;;
   stop)
	store
	;;
      *)
	echo "Usage: alsasettings stop"
	;;
esac

exit 0
