# audio/alsa-utils/init.d/alsasettings
#!/bin/bash
#
# alsasettings  Takes care of restoring the ALSA mixer values
#		2.6.x kernel systems
#
# chkconfig: 2345 98 12
# description: Restore ALSA settings on 2.6.x based systems

PATH="$PATH:/usr/sbin:/usr/bin"

function restore() {
  echo "Restoring ALSA mixer settings..."
  alsactl restore
}

function store() {
  echo "Storing ALSA mixer settings..."
  alsactl store
  amixer set Master mute &> /dev/null
}

case "$1" in
  start)
	restore
	;;
   stop)
	store
	;;
      *)
	echo "Usage: alsasettings {start|stop}"
	;;
esac

exit 0
