# zbeta/tor/init.d/tor
#!/bin/sh
#
# short: tor
# chkconfig: 2345 90 10
# processname: tor
# description: The Onion Router - A low-latency anonymous proxy
# pidfile: /var/run/tor.pid
# runas: tor
#

. /lib/lsb/init-functions

