# gnome2/apps/gnumeric/profile.d/gnumeric.rc
export PATH=$PATH:/opt/lunar/gnumeric/bin
export XDG_DATA_DIRS=$XDG_DATA_DIRS:/opt/lunar/gnumeric/share
