# utils/uptimed/systemd.d/uptimed.service
[Unit]
Description=system uptime record daemon

[Service]
ExecStart=/usr/sbin/uptimed -f

[Install]
WantedBy=multi-user.target
