# crypto/heimdal/init.d/kdc
#!/bin/bash
#
# chkconfig: 345 30 70
# description: Heimdal KDC daemon
# processname: kdc
# pidfile: /var/run/kdc.pid

ARGS="--detach"

. /lib/lsb/init-functions

