# zbeta/e17/eina-svn/profile.d/e17-cvs.rc
PATH=${PATH}:/opt/lunar/enlightenment/bin
PKG_CONFIG_PATH=${PKG_CONFIG_PATH}:/opt/lunar/enlightenment/lib/pkgconfig

export PATH PKG_CONFIG_PATH
