# libs/glibc/init.d/nscd
#!/bin/bash

# chkconfig: 2345 20 80
# description: nscd is the glibc name service caching daemon - see man nscd
# processname: /usr/sbin/nscd
# sigreload: HUP

start ()
{
	echo -n "Starting nscd: "
	/usr/sbin/nscd &&
	echo -e $RESULT_OK || echo -e $RESULT_FAIL
}

stop ()
{
	echo -n "Stopping nscd: "
	/usr/sbin/nscd --shutdown &&
	echo -e $RESULT_OK || echo -e $RESULT_FAIL
}

status ()
{
	echo -n "Service nscd: "
	/bin/pgrep -P 1 nscd > /dev/null &&
	echo -e $RESULT_OK || echo -e $RESULT_FAIL
}

. /lib/lsb/init-functions $1

