# crater/bluez/init.d/bluetoothd
#!/bin/sh
#
#   Startup/shutdown script for Bluetooth daemons
#
# chkconfig: 345 80 20
# description: Startup/shutdown script for the Bluetooth daemons
# processname: /usr/sbin/bluetoothd
# config: /etc/bluetooth/main.conf

 . /lib/lsb/init-functions $1


