# web/apache/init.d/httpd
#!/bin/bash
#
# Startup script for the Apache Web Server
#
# chkconfig: - 85 15
# description: Apache is a World Wide Web server.  It is used to serve \
#              HTML files and CGI.
# processname: httpd
# pidfile: /var/run/httpd.pid
# config: /etc/httpd/conf/access.conf
# config: /etc/httpd/conf/httpd.conf
# config: /etc/httpd/conf/srm.conf

case $1 in
     start|restart|stop)  echo      "$1ing Apache web server."
                          mkdir  -p  /var/run/httpd
                          apachectl  $1
                          ;;
                      *)  echo      "Usage: $0 {start|stop|restart}"
                          ;;
esac
