# net/mDNSResponder/init.d/mdnsd
#!/bin/sh
#
# Startup file for the mdns responder
#
# chkconfig: 345 85 15
# description: mdnsd is Apple's zeroconf responder daemon
# processname: mdnsd
# pidfile: /var/run/mdnsd.pid

. /lib/lsb/init-functions $1
