# ftp/vsftpd/systemd.d/vsftpd.service
[Unit]
Description=vsftpd ftp server
After=network.target

[Service]
ExecStart=/usr/sbin/vsftpd

[Install]
WantedBy=multi-user.target
