# zdeprecated/gnustep/gnustep-make/profile.d/gnustep.rc
export GNUSTEP_SYSTEM_ROOT=/opt/lunar/gnustep/1/
# before load GNUSTEP variables
. $GNUSTEP_SYSTEM_ROOT/System/Library/Makefiles/GNUstep.sh
export PATH=$GNUSTEP_SYSTEM_ROOT/Tools:$PATH
