# chat/dancer-ircd/init.d/dancer-ircd
#!/bin/bash
#
# irc     - the standard IRC daemon
#
# chkconfig: 345 80 20
# description: irc daemon
# processname: dancer-ircd
# config:      /etc/ircd.conf
# runas:       ircd


ARGS="-h `/bin/hostname` -f /etc/dancer-ircd/ircd.conf"

cd /

. /lib/lsb/init-functions

