# zbeta/upstart/init/tty1.conf
start on stopped rc RUNLEVEL=[23456]
stop on runlevel [!23456]

respawn
exec /sbin/agetty 38400 tty1
