# zbeta/upstart/init/tty4.conf
start on runlevel [23]
stop on runlevel [!23]

respawn
exec /sbin/agetty 38400 tty4
