# ftp/vsftpd/init.d/vsftpd
#!/bin/bash
#
# vsftpd        VSftp ftp server
#
# chkconfig: 2345 80 20
# description: ftp daemon
# processname: vsftpd
# config:      /etc/vsftpd.conf
# pidfile:     /var/run/vsftpd.pid


. /lib/lsb/init-functions

