# chat/reaim/init.d/reaim
#!/bin/bash
#
# reaim   ReAIM AIM Transparent Proxy
#
# chkconfig: 35 80 20
# description: aim proxy
# processname: reaim

if [ -e /etc/config.d/reaim ]; then
    . /etc/config.d/reaim
    ARGS=$IP
fi;

. /lib/lsb/init-functions

