# chat/reaim/init.d/reaim
#!/bin/bash
#
# reaim   ReAIM AIM Transparent Proxy
#
# chkconfig: 35 80 20
# description: aim proxy
# processname: reaim

. /lib/lsb/init-functions

