# web/apache-mod_ssl/init.d/apache_modssl
#!/bin/sh
#
# Startup script for the Apache Web Server SSL
#
# chkconfig: 345 85 15
# description: Apache is a World Wide Web server.  It is used to serve \
#              HTML files and CGI.
# processname: apache_modssl
# pidfile: /var/run/httpsd/httpd.pid
# config: /etc/httpsd/conf/httpd.conf


case $1 in
     start         )  echo      "$ing Apache web server."
                      mkdir  -p  /var/run/httpsd
                      apachectl  $1ssl
                      ;;
     stop          )  echo      "$ing Apache web server."
                      apachectl  $1
                      ;;
     restart       )  echo      "$ing Apache web server."
                      mkdir  -p  /var/run/httpsd
                      apachectl  stop
                      sleep 2
                      apachectl  startssl
                      ;;
                  *)  echo      "Usage: $0 {start|stop|restart}"
                      ;;
esac
