# chat/imspector/init.d/imspector
# imspector 
#  
# chkconfig: 345 90 10 
# description: Instant Messenger transparent proxy with logging capabilities. \
# Currently it supports MSN, AIM, ICQ, Yahoo and IRC to different degrees.
#
# pidfile: /var/run/imspector.pid
# processname: imspector
# config: /etc/imspector.conf

ARGS="-c /etc/imspector/imspector.conf"
. /lib/lsb/init-functions
