# filesys/udev/49-udev.rules
# Lunar's rules udev file
# this file assumes you have
# * root.root as default user/group
# * 0660 as default perms
# for each node

# console
KERNEL="console",                                 MODE="0600"                      
KERNEL="root[0-9][0-9]*",        OWNER="$local"
KERNEL="vc/[0-9]*",              OWNER="$local"                 
KERNEL="ptmx",                   OWNER="$local",  MODE="0666"
KERNEL="pty/m*",                 OWNER="$local" 
KERNEL="root[p-za-e][0-9a-f]*",  OWNER="$local"
KERNEL="vc/s*",                  OWNER="$local"                                 


# cpu devices
KERNEL="cpu[0-9]*",  NAME="cpu/%n/cpuid"
KERNEL="msr[0-9]*",  NAME="cpu/%n/msr"
KERNEL="microcode",  NAME="cpu/microcode"


# tty devices
KERNEL="tty",                                                         GROUP="tty", MODE="0666"
KERNEL="tty[0-9]*",   NAME="vc/%n",   SYMLINK="%k",  OWNER="$local",  GROUP="tty"
KERNEL="ttyS[0-9]*",  NAME="tts/%n",  SYMLINK="%k",  OWNER="$local",  GROUP="tty"


# pty devices
KERNEL="tty[p-za-e][0-9a-f]*",  NAME="tty/s%n",  SYMLINK="%k",  OWNER="$local" 
KERNEL="pty[p-za-e][0-9a-f]*",  NAME="pty/m%n",  SYMLINK="%k",  OWNER="$local"                             


# vc devices
KERNEL="vcc/*",                                      MODE="0600"                                
KERNEL="vcs",         NAME="vcc/0",   SYMLINK="%k",  MODE="0600"
KERNEL="vcs[0-9]*",   NAME="vcc/%n",  SYMLINK="%k",  MODE="0600"
KERNEL="vcsa",                                       MODE="0600"
KERNEL="vcsa[0-9]*",  NAME="vcc/a%n", SYMLINK="%k",  MODE="0600"


# mem
KERNEL="mem",      GROUP="kmem",  MODE="0640"
KERNEL="kmem",     GROUP="kmem",  MODE="0640"
KERNEL="port",     GROUP="kmem",  MODE="0640"
KERNEL="random",                  MODE="0666"
KERNEL="urandom",                 MODE="0444"
KERNEL="zero",                    MODE="0666"
KERNEL="null",                    MODE="0666"
KERNEL="full",                    MODE="0666"


# miscelaneous
KERNEL="agpgart",      NAME="misc/%k",  SYMLINK="%k", GROUP="video"
KERNEL="rtc",          NAME="misc/%k",  SYMLINK="%k", MODE="0664"
KERNEL="inotify",      NAME="misc/%k",  SYMLINK="%k", MODE="0666"
KERNEL="psaux",        NAME="misc/%k",  SYMLINK="%k"
KERNEL="uinput",       NAME="misc/%k",  SYMLINK="%k"
KERNEL="nvram",        NAME="misc/%k"


# floppy
KERNEL="fd[0-9]*", NAME="floppy/%n", SYMLINK="%k", GROUP="floppy"
KERNEL="floppy/*",                                 GROUP="floppy"


# audio

# alsa
KERNEL="controlC[0-9]*",  NAME="snd/%k",  GROUP="audio" 
KERNEL="hw[CD0-9]*",      NAME="snd/%k",  GROUP="audio"   
KERNEL="pcm[CD0-9cp]*",   NAME="snd/%k",  GROUP="audio"
KERNEL="timer",           NAME="snd/%k",  GROUP="audio"
KERNEL="seq",             NAME="snd/%k",  GROUP="audio"
KERNEL="midiC[D0-9]*",    NAME="snd/%k",  GROUP="audio"

# sound 
#KERNEL="admm*",                                                 GROUP="audio"
KERNEL="adsp",        NAME="sound/adsp",         SYMLINK="%k",   GROUP="audio"
KERNEL="adsp*",       NAME="sound/adsp%n",       SYMLINK="%k",   GROUP="audio"
KERNEL="aload*",                                                 GROUP="audio"
KERNEL="amidi*",                                                 GROUP="audio"
KERNEL="audio",       NAME="sound/audio",        SYMLINK="%k",   GROUP="audio"
KERNEL="audio*",      NAME="sound/audio%n",      SYMLINK="%k",   GROUP="audio"
KERNEL="beep",                                                   GROUP="audio"
KERNEL="dmfm*",                                                  GROUP="audio"
KERNEL="dmmidi*",                                                GROUP="audio"
KERNEL="dsp",         NAME="sound/dsp",          SYMLINK="%k",   GROUP="audio"
KERNEL="dsp*",        NAME="sound/dsp%n",        SYMLINK="%k",   GROUP="audio"
KERNEL="midi*",                                                  GROUP="audio"
KERNEL="mixer",       NAME="sound/mixer",        SYMLINK="%k",   GROUP="audio"
KERNEL="mixer*",      NAME="sound/mixer",        SYMLINK="%k",   GROUP="audio"
KERNEL="sequencer*",  NAME="sound/sequencer%n",  SYMLINK="%k",   GROUP="audio"
KERNEL="sequencer",   NAME="sound/sequencer",    SYMLINK="%k",   GROUP="audio"
KERNEL="snd/*",                                                  GROUP="audio"
KERNEL="sndstat",                                                GROUP="audio"
KERNEL="sound/*",                                                GROUP="audio"


# pilot / palm
KERNEL="pilot",  GROUP="uucp"


# jaz devices
KERNEL="jaz*",  GROUP="disk"


# zip devices
KERNEL="zip*",        GROUP="disk"
KERNEL="pocketzip*",  GROUP="disk"


# ls120 devices
KERNEL="ls120",   GROUP="disk"
KERNEL="ls120*",  GROUP="disk"


# lp devices
KERNEL="irlpt*",    GROUP="lp"
KERNEL="lp*",       GROUP="lp"
KERNEL="parport*",  GROUP="lp"
KERNEL="usb/lp*",   GROUP="lp"
KERNEL="usblp*",    GROUP="lp"


# scanner devices
KERNEL="scanner",                     MODE="0600"
KERNEL="usb/scanner*",  GROUP="usb",  MODE="0600"


# camera devices
KERNEL="camera*",                    MODE="0600"
KERNEL="usb/dc2xx*",   GROUP="usb",  MODE="0600"
KERNEL="usb/mdc800*",  GROUP="usb",  MODE="0600"
KERNEL="raw1394",      GROUP="usb",  MODE="0600"


# raw devices
KERNEL="raw/*",  GROUP="disk"


# ram devices
KERNEL="ram",       NAME="rd/0",   SYMLINK="%k",  GROUP="disk"
KERNEL="ram[0-9]*", NAME="rd/%n",  SYMLINK="%k",  GROUP="disk"


# disk devices
BUS="ide",  KERNEL="hd[a-z]",  SYSFS{removable}="1", PROGRAM="/etc/udev/scripts.d/ide-devfs.sh %k %b %n", NAME="%k", SYMLINK="%c{1} %c{2} cdrom%e", GROUP="cdrom"
BUS="ide",  KERNEL="hd*",                            PROGRAM="/etc/udev/scripts.d/ide-devfs.sh %k %b %n", NAME="%k", SYMLINK="%c{1} %c{2}",         GROUP="disk"

BUS="scsi", KERNEL="sr*", SYSFS{type}="5", PROGRAM="/etc/udev/scripts.d/scsi-devfs.sh sr %b %n", NAME="%c{1}", SYMLINK="%c{2} %k %c{3} %c{4} cdrom%e", GROUP="cdrom"
BUS="scsi", KERNEL="sg*", SYSFS{type}="5", PROGRAM="/etc/udev/scripts.d/scsi-devfs.sh sg %b %n", NAME="%c{1}", SYMLINK="%c{2} %k %c{3} %c{4} cdrom%e", GROUP="cdrom"
## IAC, scsi devices are disks ##
BUS="scsi", KERNEL="sd*", PROGRAM="/etc/udev/scripts.d/scsi-devfs.sh sd %b %n", NAME="%c{1}", SYMLINK="%c{2} %k %c{3} %c{4}", GROUP="disk"
BUS="scsi", KERNEL="sr*", PROGRAM="/etc/udev/scripts.d/scsi-devfs.sh sr %b %n", NAME="%c{1}", SYMLINK="%c{2} %k %c{3} %c{4}", GROUP="disk"
BUS="scsi", KERNEL="st*", PROGRAM="/etc/udev/scripts.d/scsi-devfs.sh st %b %n", NAME="%c{1}", SYMLINK="%c{2} %k %c{3} %c{4}", GROUP="disk"
BUS="scsi", KERNEL="sg*", PROGRAM="/etc/udev/scripts.d/scsi-devfs.sh sg %b %n", NAME="%c{1}", SYMLINK="%c{2} %k %c{3} %c{4}", GROUP="disk"

KERNEL="dasd*",          GROUP="disk"
KERNEL="discs/*/*",      GROUP="disk"
KERNEL="ide/*/*/*/*/*",  GROUP="disk"
KERNEL="ataraid*",       GROUP="disk"
KERNEL="pg*",            GROUP="disk"
KERNEL="microdrive*",    GROUP="disk"

# raid devices
KERNEL="cciss!*", PROGRAM="/etc/udev/scripts.d/raid-devfs.sh %k", NAME="%c{1}", SYMLINK="%k", GROUP="disk"
KERNEL="ida!*",   PROGRAM="/etc/udev/scripts.d/raid-devfs.sh %k", NAME="%c{1}", SYMLINK="%k", GROUP="disk"
KERNEL="rd!*",    PROGRAM="/etc/udev/scripts.d/raid-devfs.sh %k", NAME="%c{1}", SYMLINK="%k", GROUP="disk"

# loop devices
KERNEL="loop",        NAME="loop/0",   SYMLINK="%k",  GROUP="disk"
KERNEL="loop[0-9]*",  NAME="loop/%n",  SYMLINK="%k",  GROUP="disk"
KERNEL="loop/*",                                      GROUP="disk"

# optical devices
KERNEL="cdrom*",                                  GROUP="cdrom"
KERNEL="cdroms/*",  NAME="%k", SYMLINK="cdrom%e", GROUP="cdrom"
KERNEL="cdwriter",  NAME="%k", SYMLINK="cdrom%e", GROUP="cdrom"
KERNEL="sr*",       NAME="%k", SYMLINK="cdrom%e", GROUP="cdrom"
KERNEL="dvd",       NAME="%k", SYMLINK="cdrom%e", GROUP="cdrom"
KERNEL="rdvd",      NAME="%k", SYMLINK="cdrom%e", GROUP="cdrom"
KERNEL="pcd*",      NAME="%k", SYMLINK="cdrom%e", GROUP="cdrom"
KERNEL="scd*",      NAME="%k", SYMLINK="cdrom%e", GROUP="cdrom"


# md block devices
KERNEL="md[0-9]*",  NAME="md/%n",  SYMLINK="%k",  GROUP="disk"
KERNEL="md/*",                                    GROUP="disk"


# tape
KERNEL="nht*",      GROUP="tape"
KERNEL="npt*",      GROUP="tape"
KERNEL="nst*",      GROUP="tape"
KERNEL="pt[0-9]*",  GROUP="tape"
KERNEL="st*",       GROUP="tape"
KERNEL="ht*",       GROUP="tape"


# dm devices
KERNEL="dm-[0-9]*", PROGRAM="/sbin/devmap_name %M %m", NAME="%k", SYMLINK="%c",  MODE="0640"


# memstick
KERNEL="memstick*", MODE="0600"


# flash devices
KERNEL="flash*", MODE="0600"


# diskonkey devices
KERNEL="diskonkey*",  GROUP="disk"


# fb devices
KERNEL="fb[0-9]*", NAME="fb/%n", SYMLINK="%k", MODE="0600"
KERNEL="fb/*",                                 MODE="0600"


# v4l devices
KERNEL="v4l/*",                            GROUP="video"
KERNEL="video[0-9]*", NAME="v4l/video%n",  GROUP="video"
KERNEL="video/*",                          GROUP="video"
KERNEL="vttuner",                          GROUP="video"
KERNEL="vtx[0-9]*",   NAME="v4l/vtx%n",    GROUP="video"
KERNEL="vtx*",                             GROUP="video"
KERNEL="winradio*",                        GROUP="video"
KERNEL="radio[0-9]*", NAME="v4l/radio%n",  GROUP="video"
KERNEL="vbi[0-9]*",   NAME="v4l/vbi%n",    GROUP="video"
KERNEL="vbi*",                             GROUP="video"


# input devices
KERNEL="input/*",  MODE="0644"


# gpm devices
KERNEL="gpmctl",  MODE="0700"


# dri devices
KERNEL="3dfx*",                        GROUP="video"
KERNEL="card*",    NAME="dri/card%n",  GROUP="video" 
KERNEL="nvidia*",                      GROUP="video"
KERNEL="dri/*",                        GROUP="video"


# mainboard devices
KERNEL="apm_bios",                                  MODE="0600"
KERNEL="i2c-[0-9]*",  NAME="i2c/%n", SYMLINK="%k"


# serial + dialup devices
#KERNEL="ippp*"
#KERNEL="ircomm*"
#KERNEL="isdn*"
#KERNEL="isdnctrl*"
#KERNEL="tts/[0-9]*"
#KERNEL="capi*"
#KERNEL="dcbri*"
#KERNEL="rootS*"


# usb
KERNEL="auer*",          NAME="usb/%k",  GROUP="usb"
KERNEL="dabusb*",        NAME="usb/%k",  GROUP="usb"
KERNEL="hiddev*",        NAME="usb/%k",  GROUP="usb"
KERNEL="legousbtower*",	 NAME="usb/%k",  GROUP="usb"
KERNEL="ttyUSB[0-9]*",   NAME="usb/%k",  GROUP="usb"      
KERNEL="scanner[0-9]*",  NAME="usb/%k",  GROUP="usb"
KERNEL="tts/USB[0-9]*",                  GROUP="usb"
KERNEL="usb/*",                          GROUP="usb"


# input
KERNEL="ts*",     NAME="input/%k",                 MODE="0600"
KERNEL="event*",  NAME="input/%k",                 MODE="0600"
KERNEL="djs*",	  NAME="input/%k",                 MODE="0644"
KERNEL="js*",	  NAME="input/%k",                 MODE="0644"
KERNEL="kbd",	  NAME="input/%k",                 MODE="0644"
KERNEL="mice",	  NAME="input/%k",  SYMLINK="%k",  MODE="0644"
KERNEL="mouse*",  NAME="input/%k",                 MODE="0644"


# dvb devices
KERNEL="dvb*",  PROGRAM="/etc/udev/scripts.d/dvb.sh %k",  NAME="%c",  GROUP="video"

ENV{UDEVD_EVENT}=="1", RUN+="/sbin/udev_run_hotplugd"
