# net/samba/init.d/smbd
#!/bin/bash

# chkconfig: 345 80 20
# description: smbd is the samba file serving daemon
# processname: smbd
# config: /etc/samba/smb.conf

ARGS="-D"

. /lib/lsb/init-functions $1

