# net/dante/init.d/sockd
#!/bin/bash
#
# sockd SOCKS daemon
#
# chkconfig: 345 80 20
# description: SOCKS daemon
# processname: sockd
# config:      /etc/sockd.conf

ARGS="-D"

. /lib/lsb/init-functions

