# science/seismicunix/profile.d/seismicunix.rc
export CWPROOT="/opt/lunar/su"
PATH="${PATH}:/opt/lunar/su/bin"
