# utils/sysklogd/init.d/klogd
#!/bin/bash
#
# klogd        Starts klogd.
#
#
# chkconfig: 2345 12 88
# processname: klogd
# description: the kernel log daemon allows kernel messages to be logged \
# through syslog

. /lib/lsb/init-functions
