# utils/sysklogd/init.d/klogd
#!/bin/bash
#
# klogd        Starts klogd.
#
#
# chkconfig: 2345 12 88
# processname: klogd

. /lib/lsb/init-functions
