# utils/sysklogd/init.d/syslogd
#!/bin/bash
#
# syslogd        Starts syslogd
#
#
# chkconfig: 2345 12 88
# processname: syslogd
# description: syslogd is the facility by which many daemons use to log \
# messages to various system log files.  It is a good idea to always \
# run syslog.

. /lib/lsb/init-functions

