# shells/bash/profile.d/bash-aliases-base.rc
alias ls='ls --color=auto'
alias 'cd..'='cd ..'
alias '..'='cd ..'
alias '...'='cd ../..'
alias '/'='cd /'
alias '~'='cd'
