# utils/smartmontools/init.d/smartd
#! /bin/sh
#
# chkconfig: 2345 40 40
# description: Self Monitoring and Reporting Technology (SMART) Daemon
# processname: /usr/sbin/smartd 
#
. /lib/lsb/init-functions 
