# utils/powertweak/init.d/powertweakd
#!/bin/bash
#
# powertweak - configure power stuff or something
#
# processname: powertweakd
# chkconfig: 345 70 30

. /lib/lsb/init-functions

