# kde4/base/kdelibs/init.d/kdm4
#!/bin/bash

# chkconfig: 5 95 05
# description: kdm is KDE's version of the xdm graphical login manager
# processname: /usr/bin/kdm
# sigreload: HUP

. /lib/lsb/init-functions $1
