# gnome2-apps/gnumeric2/profile.d/gnumeric.rc
export PATH=$PATH:/opt/lunar/gnumeric/bin
export PKG_CONFIG_PATH:$PKG_CONFIG_PATH:/opt/lunar/gnumeric/lib/pkgconfig
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/opt/lunar/gnumeric/lib
