# net/rabbitmq-server/init.d/rabbitmq-server
#!/bin/bash
#
# rabbitmq-server RabbitMQ broker
#
# chkconfig: - 80 05
# description: Enable AMQP service provided by RabbitMQ
# processname: rabbitmq-multi

DAEMON=/usr/sbin/rabbitmq-multi
INIT_LOG_DIR=/var/log/rabbitmq
NODE_COUNT=1

start() {
  echo -n "Starting rabbitmq-server: "
  $DAEMON start_all ${NODE_COUNT} > ${INIT_LOG_DIR}/startup_log 2> ${INIT_LOG_DIR}/startup_err &&
  echo -e "$RESULT_OK" || echo -e "$RESULT_FAIL"
}

stop() {
  status_rabbitmq quiet
  if [ $? = 0 ]; then
    echo -n "Stopping rabbitmq-server: "
    $DAEMON stop_all > ${INIT_LOG_DIR}/shutdown_log 2> ${INIT_LOG_DIR}/shutdown_err &&
    echo -e "$RESULT_OK" || echo -e "$RESULT_FAIL"
  else
    echo "No nodes running"
    exit 0
  fi
}

status_rabbitmq() {
    if [ "$1" != "quiet" ] ; then
        $DAEMON status 2>&1
    else
        $DAEMON status > /dev/null 2>&1
    fi
}

status() {
  status_rabbitmq
}

rotate_logs() {
  echo -n "Rotating rabbitmq-server logs: "
  $DAEMON rotate_logs &> /dev/null &&
  echo -e "$RESULT_OK" || echo -e "$RESULT_FAIL"
}

restart() {
  $0 stop
  sleep 1
  $0 start
}

case "$1" in
  start|stop|restart|rotate_logs|status) ;;
  *) echo "Usage: $0 {start|stop|restart|rotate_logs|status}" ; exit 1;;
esac

. /lib/lsb/init-functions
