# mail/imap/xinetd.d/imap
service imap2
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/sbin/imapd
	per_source	= 8
	log_type	= FILE /var/log/xinetd/imap
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
