# utils/metalog/init.d/syslog
#!/bin/bash
#
# metalog        Starts the logger.
#
#
# chkconfig: 2345 12 88
# description: metalog is the facility by which many daemons use to log \
# messages to various system log files.  It is a good idea to always \
# run metalog.

case $1 in
  start|restart)  echo     "$1ing metalog daemons."
                  pkill -f  "metalog \[M"
                  metalog  --daemonize  --synchronous
                  ;;

   stop)          echo     "$1ping metalog daemons."
                  pkill -f  "metalog \[M"
                  ;;                 

      *)          echo     "Usage: $0 {start|stop|restart}"
                  ;;                                        
esac
