# net/netatalk/init.d/netatalk
#!/bin/bash
#
# chkconfig: 345 85 15
# description: AFP server
# processname: netatalk
# config: /etc/afp.conf

. /lib/lsb/init-functions
