# utils/gpm/profile.d/gpm.rc
case $(/usr/bin/tty) in
  /dev/tty[0-9]*) [ -n "$(pidof -s gpm &>/dev/null)" ] && /usr/bin/disable-paste ;;
esac
