# web/webmin/init.d/webmin
#!/bin/sh
#
# Start/stop Webmin
#
# chkconfig: 345 90 10
# description: Webmin WEB Admin Tool
# processname: /var/webmin/html/miniserv.pl
# pidfile: /var/webmin/miniserv.pid

ARGS="/etc/webmin/miniserv.conf"

restart () {
  $0 stop
  $0 start
}

. /lib/lsb/init-functions
 
