# java/clojurescript/profile.d/clojurescript.rc
export CLOJURESCRIPT_HOME="/usr/share/java/clojurescript"
export PATH+=":/usr/share/java/clojurescript/bin"
