# java/clojurescript/profile.d/clojurescript.rc
export CLOJURESCRIPT_HOME=/usr/share/java/clojurescript
export CLASSPATH=$CLASSPATH:$CLOJURESCRIPT_HOME/lib/clojurescript.jar