# java/clojurescript/profile.d/clojurescript.rc
export CLOJURESCRIPT_HOME="/usr/share/java/clojurescript"
export CLASSPATH="$CLOJURESCRIPT_HOME/lib/clojurescript.jar:$CLASSPATH"
export PATH+=":/usr/share/java/clojurescript/bin"
