# utils/lard/init.d/lard
#! /bin/bash
#
# lard          Start/Stop the lard logging daemon.
#
# chkconfig: 2345 12 88
# description: Lard is a Logging and Rotation Daemon
# processname: lard
# pidfile: /var/run/lard.pid

. /lib/lsb/init-functions
