# utils/fcron/init.d/crond
#! /bin/bash
#
# fcron          Start/Stop the cron clock daemon.
#
# chkconfig: 2345 90 60
# description: cron is a standard UNIX program that runs user-specified \
#              programs at periodic scheduled times. vixie cron adds a \
#              number of features to the basic UNIX cron, including better \
#              security and more powerful configuration options.
# processname: /usr/sbin/fcron
# pidfile: /var/run/fcron.pid

ARGS="-b"

. /lib/lsb/init-functions
