# utils/smartmontools/systemd.d/smartd.service
[Unit]
Description=smartd control and monitor your harddiscs

[Service]
ExecStart=/usr/sbin/smartd -n

[Install]
WantedBy=multi-user.target
