# crater/kotlin/profile.d/kotlin.rc
#!/usr/bin/env sh
export KOTLIN_HOME=/usr/share/java/kotlin
