# net/arpwatch/init.d/arpwatch
#!/bin/bash
#
# arpwatch	Watches mac addrs deamon
#
# chkconfig: 345 80 20
# description: arpwarch daemon
# processname: arpwatch
#

ARGS="-f /var/state/arpwatch/arp.dat"

. /lib/lsb/init-functions

