# net/rpcbind/init.d/rpcbind
#!/bin/sh
#
# rpcbind
#
# chkconfig: 345 11 89
# description: the rpcbind is a RPC port broker
# processname: /usr/bin/rpcbind
#

. /lib/lsb/init-functions

