# utils/colorgcc/profile.d/colorgcc.rc
# prefix PATH with colorgcc overlay
export PATH="/usr/lib/colorgcc:$PATH"
