# bluetooth/bluez-utils/init.d/btdund
#!/bin/bash
#
# btdund Bluetooth DUND daemon
#
# chkconfig: 345 80 20
# description: SDP Bluetooth Daemon
# processname: dund

ARGS="--listen -s --channel 1"

. /lib/lsb/init-functions

