# bluetooth/bluez-utils/init.d/bthcid
#!/bin/bash
#
# bthcid Bluetooth HCI Daemon
#
# chkconfig: 345 80 20
# description: HCI Bluetooth Daemon
# processname: hcid
# config:      /etc/bluetooth/hcid.conf

ARGS="-f /etc/bluetooth/hcid.conf"

. /lib/lsb/init-functions

