# bluetooth/bluez-utils/init.d/btrfcomm
#!/bin/bash
#
# rfcomm
#
# chkconfig: 345 80 20
# description: rfcomm Bluetooth Daemon
# processname: rfcomm
# config:      /etc/bluetooth/rfcomm.conf

ARGS="-f /etc/bluetooth/rfcomm.conf"

. /lib/lsb/init-functions

