# mail/amavisd-new/init.d/amavisd
#!/bin/bash
#
# amavisd
# chkconfig: 2345 78 32
# description: amavisd is the amavis daemon
# processname: amavisd
# config: /etc/amavis.conf
# pidfile: /var/amavis/amavisd.pid

. /lib/lsb/init-functions

