# gnome-desktop/gdm/init.d/gdm
#!/bin/bash

# chkconfig: 5 95 05
# description: gdm is gnome's version of the xdm graphical login manager
# processname: /opt/lunar/gnome/2/bin/gdm-binary
# sigreload: HUP

. /lib/lsb/init-functions $1
