# distributed/distcc/init.d/distcc
#!/bin/bash
#
# distcc  -  the distributed compiling daemon
#
# chkconfig: 345 80 20
# description: distcc daemon
# processname: distccd
# runas:       nobody

reload () {
	default_restart
}

. /lib/lsb/init-functions

