# distributed/distcc/init.d/distcc
#!/bin/bash
#
# distcc  -  the distributed compiling daemon
#
# chkconfig: 345 80 20
# description: distcc daemon
# processname: distccd

ARGS="--daemon --user nobody --nice 19 -P /var/run/distccd.pid --allow 127.0.0.1"

# Append DISTCC_ALLOW values from distcc optimizations to daemon args
if [ -f /etc/lunar/local/optimizations.DISTCC ]; then
  . /etc/lunar/local/optimizations.DISTCC

  for host_or_ip_range in $DISTCC_ALLOW; do
    ARGS="$ARGS --allow $host_or_ip_range"
  done
fi

reload () {
	default_restart
}

. /lib/lsb/init-functions

