# distributed/distcc/init.d/distcc
#!/bin/bash
#
# distcc  -  the distributed compiling daemon
#
# chkconfig: 345 80 20
# description: distcc daemon
# processname: distccd

ARGS="--daemon --user nobody --nice 19 -P /var/run/distccd.pid"

reload () {
	default_restart
}

. /lib/lsb/init-functions

