# zbeta/xmms2-beta/profile.d/xmms2.rc
# Ugly hack to start xmms2d for the user that logon
if ! pidof xmms2d 2>&1 > /dev/null;  then
  /usr/bin/xmms2-launcher
fi
