# net/portmap/systemd.d/portmap.service
[Unit]
Description=RPC Port Mapper
After=network.target

[Service]
ExecStart=/usr/sbin/portmap -d

[Install]
WantedBy=multi-user.target
