# devel/hal/init.d/hald
#!/bin/bash
#
# Startup script for hald
#
# description: Hardware Abstraction Layer
# chkconfig: 345 98 02
# processname: hald
# pidfile: /var/run/hald.pid

ARGS="--daemon=yes --use-syslog"

. /lib/lsb/init-functions $1
