# audio/mpd/init.d/mpd
#!/bin/bash
#
# Startup script for mpd
#
# chkconfig: 345 99 10
# short: musicmpd
# description: MusicPD daemon
# processname: mpd

stop() {
	echo -n "Stopping $SHORT:"
	mpd --kill &> /dev/null
	echo -e $RESULT_OK
}

. /lib/lsb/init-functions $1
