# mail/sendmail/init.d/sendmail
#!/bin/bash
#
# sendmail	sendmail SMTP server
#
# chkconfig: 2345 80 20
# description: sendmail is the standard SMTP server for sending and \
#              receiving mail across the internet.
# processname: sendmail
# pidfile:     /var/run/sendmail.pid
# config:      /etc/mail/sendmail.cf

ARGS="-bd -q5m"

reload () {
	# not reload safe!
	default_restart
}

start () {
	if [ -f /etc/mail/aliases ] ; then
	  newaliases
	fi
	if [ -f /etc/mail/virtusers ] ; then
	  makemap hash /etc/mail/virtusers.db < /etc/mail/virtusers
	fi
	if [ -f /etc/mail/genericstable ] ; then
	  makemap hash /etc/mail/genericstable.db < /etc/mail/genericstable
	fi
	if [ -f /etc/mail/access ] ; then
	  makemap hash /etc/mail/access.db < /etc/mail/access
	fi
	echo -n "Starting sendmail: "
	if [ -f $PIDFILE ] ; then
	  echo -e "$PIDFILE exists!$RESULT_FAIL"
	  exit 1
	fi
	if $(/bin/pgrep -P 1 sendmail &> /dev/null) ; then 
	  echo -e "Already running$RESULT_FAIL"
	  exit 1
	fi
	/usr/sbin/sendmail $ARGS &&
	  echo -e $RESULT_OK ||
	  echo -e $RESULT_FAIL
}

stop () {
	echo -n "Stopping sendmail: "
	if [ -f $PIDFILE ] ; then
	  PID=$(head -1 $PIDFILE)
	  if [ -d /proc/$PID ] ; then
	    if kill -TERM $PID ; then
	      rm $PIDFILE
	      echo -e $RESULT_OK
	    else
	      echo -e "kill failed$RESULT_FAIL"
	    fi
	  else
	    echo -e "can't find PID$RESULT_FAIL"
	  fi
	else
	  echo -e "can't find PID$RESULT_FAIL"
	fi
}

. /lib/lsb/init-functions

