# utils/cpufreqd/init.d/cpufreqd
#!/bin/bash
#
# cpufreqd         cpu speed daemon 
#
# chkconfig: 12345 90 10
# description: cpufreqd is a cpu throttling daemon
# processname: /usr/sbin/cpufreqd

start () {
  default_start
}


. /lib/lsb/init-functions

