# java/groovy-bin/profile.d/groovy.rc
export GROOVY_HOME=/opt/lunar/groovy
export PATH=${PATH}:${GROOVY_HOME}/bin
