# kde3/arts/profile.d/kde3.rc
export  KDEDIR=${KDE_PREFIX:-/opt/lunar/kde/3}
export  PATH=$KDEDIR/bin:$PATH
