# utils/pcsc-lite/99-pcscd-hotplug.rules
ACTION=="add", ENV{PCSCD}=="1", GROUP="pcscd", RUN+="/usr/bin/env systemctl start pcscd"
