# gnome2/apps/gnumeric2/profile.d/gnumeric.rc
export PATH=$PATH:/opt/lunar/gnumeric/bin
