# net/oidentd/init.d/oidentd
#!/bin/bash
#
# oidentd        This starts and stops oidentd
#
# chkconfig: 345 99 01
# description: oidentd is an implimentation of the auth protocol
#
# processname: /usr/sbin/oidentd

case $1 in
  start|restart)  echo   "$1ing oidentd."
                  pkill -f "oidentd -u daemon -g daemon"  &&  sleep  3
                  oidentd -u daemon -g daemon
                  ;;
           stop)  echo   "$1ping oidentd."
                  pkill -f "oidentd -u daemon -g daemon"
                  ;;
              *)  echo   "Usage: $0 {start|stop|restart}"
                  ;;
esac
