# net/avahi/systemd.d/avahi-daemon.service
[Unit]
Description=Avahi network service discovery

[Service]
ExecStart=/usr/sbin/avahi-daemon -s

[Install]
WantedBy=multi-user.target
