# mail/qpopper3/xinetd.d/qpopper
service pop3
{
       socket_type = stream  
       protocol = tcp  
       wait = no  
       user = root  
       server = /usr/sbin/popper  
       server_args = qpopper -s  
       port = 110  
} 
