# java/closure-compiler/profile.d/closure.rc
export CLOSURE_HOME="/usr/share/java/closure"
export CLASSPATH="$CLOSURE_HOME/closure.jar:$CLASSPATH"
