# net/radvd/init.d/radvd
#!/bin/bash
#
# radvd
#
# chkconfig: 345 30 70
# description: radvd is an ipv6 router advertising daemon
# processname: radvd

. /lib/lsb/init-functions

