# zdeprecated/mt-daapd/init.d/mt-daapd
#!/bin/bash

# chkconfig: 345 99 01
# description: itunesserver
# processname: mt-daapd
# config: /etc/mt-daapd.conf

ARGS=""

. /lib/lsb/init-functions $1

