# system/lunar-init/netmount
#!/bin/bash
#
#   netmount - mounts network shares on nfs and smbfs
#
# description: mount and umount the filesystem
# chkconfig: 345 20 85


start() {
  mount -t nfs,nfs4,smbfs,cifs -a && 
    echo -e "Mounting network shares:$RESULT_OK" ||
    echo -e "Mounting network shares:$RESULT_FAIL"
}

stop() {
  umount -t nfs,nfs4,smbfs,cifs -a -f -l &&
    echo -e "Unmounting network shares:$RESULT_OK" ||
    echo -e "Unmounting network shares:$RESULT_FAIL"
}

# to avoid confusion we force only these options as being valid:
case "$1" in
  start|stop) ;;
           *) echo  "Usage: $0 {start|stop}"  ;;
esac

. /lib/lsb/init-functions

