# net/openslp/systemd.d/openslp.service
[Unit]
Description=Service Location Protocol
After=network.target

[Service]
ExecStart=/usr/sbin/slpd -l /var/log/slpd.log -d
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=network.target
