# utils/clamav/init.d/freshclam
#!/bin/bash
#
# amavisd
# chkconfig: 2345 99 10
# description: freshclam is the ClamAntivirus definition auto-updater
# processname: freshclam
# config: /etc/clamd.conf

ARGS="-d -c 1"

. /lib/lsb/init-functions

