# utils/sdscsyslogd/init.d/syslog
#!/bin/bash
#
# syslogd        Starts the logger.
#
#
# chkconfig: 2345 12 88
# description: syslogd is the facility by which many daemons use to log \
# messages to various system log files.  It is a good idea to always \
# run syslogd.

case $1 in
  start|restart)  echo     "$1ing syslogd daemons."
                  pkill -f  "syslogd \[M"
                  syslogd  
                  ;;

   stop)          echo     "$1ping syslogd daemons."
                  pkill -f  "syslogd \[M"
                  ;;                 

      *)          echo     "Usage: $0 {start|stop|restart}"
                  ;;                                        
esac
