# zlocal/openoffice-bin-de/profile.d/openoffice-bin.rc
# openoffice-bin specific stuff

[ -d /opt/lunar/openoffice-de/program ] &&
    PATH="$PATH:/opt/lunar/openoffice-de/program"

