# web/apache2/init.d/httpsd
#!/bin/bash
#
# Startup script for the Apache Web Server
#
# chkconfig: 345 85 15
# description: Apache is a World Wide Web server.  It is used to serve \
#              HTML files and CGI.
# processname: httpsd
# pidfile: /var/run/httpsd.pid
# config: /etc/httpd/conf/ssl.conf

ARGS="-f /etc/httpd/conf/ssl.conf -D SSL"

. /lib/lsb/init-functions $1

