# video/DirectFB/profile.d/directfb.rc
# sets the default framebuffer device
export FRAMEBUFFER=/dev/fb0
