# gnome-platform/hal/init.d/hald
#!/bin/bash
#
# Startup script for hald
#
# description: Hardware Abstraction Layer
# chkconfig: 345 98 02
# processname: hald

ARGS=""

. /lib/lsb/init-functions $1
