# devel/pkgconfig/profile.d/pkgconfig.rc
export PKG_CONFIG_PATH=$PKG_CONFIG_PATH:/usr/lib/pkgconfig:/usr/X11R6/lib/pkgconfig
