# devel/pkgconfig/profile.d/pkgconfig.rc
export PKG_CONFIG=/usr/bin/pkg-config
export PKG_CONFIG_PATH=/opt/lunar/qt/4/lib/pkgconfig:/opt/lunar/kde/4/lib/pkgconfig:/usr/share/pkgconfig:/usr/lib/pkgconfig:$PKG_CONFIG_PATH
