# graphics/autopano-sift/profile.d/autopano-sift.rc
export PATH=$PATH:/opt/lunar/autopano-sift/bin
