# gnustep/gnustep-make/profile.d/gnustep.rc
export GNUSTEP_ROOT=/opt/lunar/gnustep/1/
# before load GNUSTEP variables
. $GNUSTEP_ROOT/System/Makefiles/GNUstep.sh
export PATH=$GNUSTEP_SYSTEM_ROOT/Tools:$PATH
