# web/usermin/init.d/usermin
#!/bin/sh
#
# Start/stop Usermin
#
# chkconfig: 345 90 10
# description: Usermin WEB Admin Tool
# processname: /var/usermin/html/miniserv.pl
# pidfile: /var/log/usermin/miniserv.pid

ARGS="/etc/usermin/miniserv.conf"

restart () {
  $0 stop
  $0 start
}

. /lib/lsb/init-functions
