# web/linuxconf/xinetd.d/linuxconf-web
service linuxconf
{
	socket_type	= stream
	protocol	= tcp
	wait		= yes
	user		= root
	server		= /bin/linuxconf
	server_args	= --http
	flags           = NOLIBWRAP
	log_type	= FILE /var/log/xinetd/linuxconf
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
