# security/clamav/init.d/clamd
#!/bin/bash
#
# amavisd
# chkconfig: 2345 75 35
# description: clamd is the ClamAntivirus Daemon
# processname: clamd
# config: /etc/clamd.conf
# pidfile: /var/run/clamd.pid

. /lib/lsb/init-functions

