# printer/lprng/init.d/lpd
#!/bin/bash

# chkconfig: 2345 60 40
# description: lpd is the print daemon required for lpr to work properly. \
#   It is basically a server that arbitrates print jobs to printer(s).
# processname: lpd
# config: /etc/printcap
# config: /etc/lpd.conf
# config: /etc/lpd.perms

. /lib/lsb/init-functions $1

