# net/connman/init.d/connmand
#! /bin/bash
#
# connmand      Start/Stop the connmand network connection manager
#
# chkconfig: 345 15 90
# description: network connection manager daemon
# processname: connmand
# pidfile: /var/run/connmand.pid

. /lib/lsb/init-functions
