# utils/vixie-cron/init.d/crond
#! /bin/bash
#
# hccron          Start/Stop the cron clock daemon.
#
# chkconfig: 2345 90 10
# description: cron is a standard UNIX program that runs user-specified \
#              programs at periodic scheduled times. vixie cron adds a \
#              number of features to the basic UNIX cron, including better \
#              security and more powerful configuration options.
# processname: crond

. /lib/lsb/init-functions

