# utils/less/profile.d/less.rc

# these gadgets make less work nicer:
export LESS="-M -S -I"
export LESSOPEN="|lesspipe.sh %s"

