# utils/less/profile.d/less.rc

# these gadgets make less work nicer:
export LESSCHARSET=latin1
export LESS="-M -S -I"
export LESSOPEN="|lesspipe.sh %s"

