# mail/spamass-milter/init.d/spamass-milter
#!/bin/bash
#
# spamass-milter
# chkconfig: 2345 77 33
# description: spamass-milter is a daemon which hooks into sendmail and routes
#              email messages to spamassassin
# processname: spamass-milter
# sigstop:     KILL

ARGS="-p /var/safe/spamass/spamass.sock -f"

reload () {
        # not reload safe!
        default_restart
}

. /lib/lsb/init-functions

