# zbeta/akka-bin/profile.d/akka.rc
export AKKA_HOME="/opt/akka"
export PATH="$PATH:$AKKA_HOME/bin"
