# x11-apps/libreoffice-bin/profile.d/libreoffice-bin.rc
# libreoffice-bin specific stuff

[ -d /opt/lunar/libreoffice/libreoffice/program/ ] &&
    PATH="$PATH:/opt/lunar/libreoffice/libreoffice/program/"

