# x11-apps/libreoffice-bin/profile.d/libreoffice-bin.rc
# libreoffice-bin specific stuff

[ -d /usr/lib/libreoffice${VERSION%.*}/program/ ] &&
    PATH="$PATH:/usr/lib/libreoffice${VERSION%.*}/program/"
