# x11-apps/libreoffice-bin/profile.d/libreoffice-bin.rc
# libreoffice-bin specific stuff

[ -d /opt/libreoffice${VERSION%.*}/program/ ] &&
    PATH="$PATH:/opt/libreoffice${VERSION%.*}/program/"
