# games/fortune/profile.d/fortune.rc
# this can be an annoyance!
/usr/games/fortune 1>&2
