# games/fortune/profile.d/fortune.rc
# this can be an annoyance!
if shopt -q login_shell; then
  /usr/bin/fortune 1>&2
fi
