# games/fortune/profile.d/fortune.rc
# this can be an annoyance!
if tty -s; then
  /usr/games/fortune 1>&2
fi
