# utils/ospmd/init.d/ospmd
#!/bin/sh
#
#   Startup/shutdown script for the acpi daemon
#
# chkconfig: 345 50 15
# description: Startup/shutdown script for the acpi daemon
# processname: /sbin/ospmd

. /lib/lsb/init-functions

