# science/gpsd/systemd.d/gpsd.service
[Unit]
Description=GPS (Global Positioning System) Daemon

[Service]
Type=simple
ExecStart=/usr/bin/gpsd -N

[Install]
Also=gpsd.socket
