# net/privoxy/init.d/privoxy
#!/bin/bash
#
# privoxy	Privoxy Proxy Server
#
# chkconfig:   2345 91 09
#
# description: Web proxy with advanced filtering capabilities \
#              such as filtering web page content, managing \
#              cookies and removing ads
#
# processname: privoxy
# config:      /etc/privoxy/config

ARGS="--pidfile /var/run/privoxy.pid --user privoxy /etc/privoxy/config"

. /lib/lsb/init-functions
