# audio/nas/systemd.d/nasd.service
[Unit]
Description=Network Audio System
After=sound.target

[Service]
ExecStart=/usr/bin/nasd

[Install]
WantedBy=multi-user.target
