# libs/cups/init.d/cupsd
#!/bin/sh
#
#   Startup/shutdown script for the Common UNIX Printing System (CUPS).
#
# chkconfig: 345 90 10
# description: Startup/shutdown script for the Common UNIX \
#              Printing System (CUPS).
# processname: /usr/sbin/cupsd

. /lib/lsb/init-functions
