# net/dhcp/init.d/dhcpd
#!/bin/bash

#
# chkconfig:    345 80 20
# description:  DHCP Server
# processname:  dhcpd
# config:       /etc/config.d/dhcpd
# config:       /etc/dhcpd.conf
#
#$DHCPD_INTERFACES=eth0
# Get service config
[ -f /etc/config.d/dhcpd ] && . /etc/config.d/dhcpd

start() {
        echo -n "Starting dhcpd: "
        /usr/sbin/dhcpd $DHCPD_ARGS && echo -e $RESULT_OK || echo -e $RESULT_FAIL
}

reload () {
        # dhcpd cannot reload
        $0 restart
}

. /lib/lsb/init-functions
