# utils/microcode_ctl/init.d/microcode
#!/bin/bash
#
# Update intel microcode
#
# chkconfig: 2345 20 80
# description: Updates the intel processor microcode
# processname: /usr/sbin/microcode_ctl

ARGS="-Qu"

start()
{
	echo -n "Loading Intel microcode: "

	# Load microcode kernel module, if needed
	grep -qo ' microcode$' /proc/misc || modprobe microcode >& /dev/null

	[ -e /dev/cpu/microcode ] &&
	/usr/sbin/microcode_ctl -Qu &&
	echo -e $RESULT_OK ||
	echo -e $RESULT_FAIL
}

stop()
{
	:
}

. /lib/lsb/init-functions

