# kernel/virtualbox-module/systemd.d/virtualbox.service
[Unit]
Description=VirtualBox guest services

[Service]
ExecStart=/usr/bin/VBox.sh

[Install]
WantedBy=basic.target
