# zbeta/knock/init.d/knockd
#!/bin/bash
#
# knockd	Port-knock Daemon
#
# chkconfig: 345 85 15
# description: port-knock daemon
# processname: knockd
# config:      /etc/konckd.conf
ARGS="-d"

. /lib/lsb/init-functions

