# mail/courier-imap/init.d/imap
#!/bin/sh
#
# This is an /etc/init.d file for Courier-IMAP
#
# chkconfig: 345 80 20
# description: Courier-IMAP - IMAP server
#

prefix=/usr
exec_prefix=/usr
bindir=${exec_prefix}/bin
sbindir=${exec_prefix}/sbin

case "$1" in
start)
        cd /
	. /etc/courier/imapd
	mkdir -p /var/lock/subsys
	touch /var/lock/subsys/courier-imap

	echo -n "Starting imap:"

	case x$IMAPDSTART in
	x[yY]*)
		# Start daemons.
		/usr/libexec/imapd.rc start
		echo -n " imap"
		;;
	esac

	. /etc/courier/imapd-ssl
	case x$IMAPDSSLSTART in
	x[yY]*)
		if test -x $COURIERTLS
		then

		# First time we start this, generate a dummy SSL certificate.

			if test ! -f $TLS_CERTFILE
			then
				echo -n " generating-SSL-certificate..."
				/etc/courier/cert/mkimapdcert >/dev/null 2>&1
			fi
			/usr/libexec/imapd-ssl.rc start
			echo -n " imap-ssl"
		fi
		;;
	esac

	POP3DSTART=""
	POP3DSSLSTART=""

	if test -f /etc/courier/pop3d
	then
		. /etc/courier/pop3d
	fi

	case x$POP3DSTART in
	x[yY]*)
		# Start daemons.
		/usr/libexec/pop3d.rc start
		echo -n " pop3"
		;;
	esac

	if test -f /etc/courier/pop3d-ssl
	then
		. /etc/courier/pop3d-ssl
	fi

	case x$POP3DSSLSTART in
	x[yY]*)
		if test -x $COURIERTLS
		then

		# First time we start this, generate a dummy SSL certificate.

			if test ! -f $TLS_CERTFILE
			then
				echo -n " generating-SSL-certificate..."
				/etc/courier/cert/mkpop3dcert >/dev/null 2>&1
			fi
			/usr/libexec/pop3d-ssl.rc start
			echo -n " pop3-ssl"
		fi
		;;
	esac

	echo ""
	;;
stop)
        echo -n "Stopping imap:"
	. /etc/courier/imapd
	. /etc/courier/imapd-ssl
	/usr/libexec/imapd.rc stop
	echo -n " imap"
	if test -x $COURIERTLS
	then
		/usr/libexec/imapd-ssl.rc stop
		echo -n " imap-ssl"
	fi

	if test -f /etc/courier/pop3d
	then
		/usr/libexec/pop3d.rc stop
		echo -n " pop3"

		if test -x $COURIERTLS
		then
			/usr/libexec/pop3d-ssl.rc stop
			echo -n " pop3-ssl"
		fi
	fi

	echo ""
	;;
restart)
	$0 stop
	$0 start
        ;;
esac
exit 0
