# net/bind/init.d/named
#!/bin/bash
#
# named           This shell script takes care of starting and stopping
#                 named (BIND DNS server).
#
# chkconfig: 345 45 55
# processname: named
# description: named (BIND) is a Domain Name Server (DNS) \
# that is used to resolve host names to IP addresses.
# probe: true
# config: /var/named/named.conf

ARGS="-u bind"

start() {
  pgrep -P 1 named || rm $PIDFILE &> /dev/null
  echo -n "Starting named:"
  /usr/sbin/named $ARGS && echo -e $RESULT_OK || echo -e $RESULT_FAIL
}

. /lib/lsb/init-functions

