# net/avahi/init.d/avahi-daemon
#!/bin/bash

# chkconfig: 345 98 02
# description: Avahi Daemon
# processname: avahi-daemon
# config: /etc/avahi/avahi-daemon.conf

ARGS="-D"

. /lib/lsb/init-functions $1

