# x11-apps/openoffice-src/profile.d/openoffice-src.rc
# openoffice-src specific stuff

[ -d /opt/lunar/openoffice-src/program ] &&
    PATH="$PATH:/opt/lunar/openoffice-src/program"
