# utils/acpid/systemd.d/acpid.service
[Unit]
Description=ACPI event daemon

[Service]
ExecStart=/usr/sbin/acpid -f

[Install]
WantedBy=multi-user.target
