# bluetooth/bluez-sdp/init.d/btsdpd
#!/bin/bash
#
# btsdpd Bluetooth SDP Daemon
#
# chkconfig: 345 80 20
# description: SDP Bluetooth Daemon
# processname: sdpd

ARGS=""

. /lib/lsb/init-functions

