# audio/ladspa_sdk/profile.d/ladspa.rc
# So "listplugins" works
export LADSPA_PATH=/usr/lib/ladspa
