# java/grails-bin/profile.d/grails.rc
export GRAILS_HOME=/opt/lunar/grails
export PATH=${PATH}:${GRAILS_HOME}/bin
