# ftp/wu-ftpd/init.d/ftpd
#!/bin/bash
#
# ftp	wu-ftp ftp server
#
# chkconfig: 2345 80 20
# description: ftp daemon
# processname: in.ftpd
# config:      /etc/ftpaccess
# config:      /etc/ftpconversions
# config:      /etc/ftpusers

ARGS="-S"

extended-status () {
	ftpcount
	ftpwho
}

usage () {
	echo "Usage: $0 {start|stop|reload|probe|restart|status|extended-status}"
}

. /lib/lsb/init-functions

