# utils/acpid/init.d/acpid
#!/bin/bash
#
# acpid        Starts acpid
#
#
# chkconfig: 2345 80 20
# processname: acpid
# description: ACPI event daemon

. /lib/lsb/init-functions

