# i18n/Canna/init.d/canna
#!/bin/sh
#
# canna         Starts the Canna kanji conversion server.
#
#
# chkconfig: 345 88 12
# processname: /usr/sbin/cannaserver
# description: Canna is the Japanese kanji conversion server.

start()
{
	# I really hate doing this
	if [ -z `pidof cannaserver` ]
	then
		rm -rf /tmp/.iroha_unix
	fi
	/usr/sbin/cannaserver
}

stop()
{
	killall /usr/sbin/cannaserver
	# And because it doesn't always clean up after itself...
	rm -rf /tmp/.iroha_unix
}

. /lib/lsb/init-functions
