# mail/qpopper/init.d/qpopper
#!/bin/bash

# chkconfig: 345 80 20
# description: Qpopper pop3 daemon
# processname: popper
# pidfile:     /var/run/popper.pid

grep -q "OPENSSL=y" /etc/lunar/local/depends/qpopper && ARGS="995 -l1 -f /etc/config.d/qpopper.conf"

reload () {
	# popper doesn't like reloads
        $0 restart
}

. /lib/lsb/init-functions

