# chat/Ultimate/init.d/ircd
#!/bin/sh
#
# ircd - Ultimate IRCd
#
# chkconfig: 345 80 20
# description: Ultimate IRC daemon
# processname: ircd
# config: /etc/ultimateircd/ircd.conf
# runas: ultimate
# pidfile: /var/run/ircd.pid

# Change directory before starting or entropy will fail
cd /var/lib/ultimateircd

. /lib/lsb/init-functions
