# mail/drac/init.d/dracd
#!/bin/sh
#
# dracd   DRAC is a daemon that dynamically updates a relay 
#         authorization map for SMTP program
#
# chkconfig: 345 80 30
# description: DRAC is a daemon that dynamically updates a relay
#              authorization map for SMTP program

case $1 in
  start|restart)  echo  -n  "$1ing rpc.dracd daemon: "
                  pkill  "^rpc.dracd$"  &&  sleep  5
                  rpc.dracd &
                  ;;

           stop)  echo  -n  "$1ping rpc.dracd daemon: "
                  pkill "^rpc.dracd$"
                  ;;

              *)  echo  "Usage: dracd {start|stop|restart}"
                  exit  1
                  ;;
esac
