# net/cfengine/init.d/cfservd
#!/bin/sh
#
#  Startup/shutdown script for the cfengine daemon
#
# chkconfig: 345 95 10
# description: cfengine server daemon
# processname: cfservd

. /lib/lsb/init-functions $1
