# kde3/qt3/profile.d/qt3.rc
export  QTDIR=/opt/lunar/qt/3
export  PATH=/opt/lunar/qt/3/bin:$PATH
