# kde3/qt3/profile.d/qt3.rc
export  QTDIR=${QT_PREFIX:-/opt/lunar/qt/3}
export  PATH=QTDIR/bin:$PATH
