# java/closure/profile.d/closure.rc
export CLOSURE_HOME=/usr/share/java/closure
export CLASSPATH=$CLASSPATH:$CLOSURE_HOME/closure.jar
