# zbeta/powernowd/init.d/powernowd
#!/bin/bash
#
# powernowd        CPUFreq client daemon
#
# chkconfig: 2345 90 10
# description: powernowd is a daemon to control the speed of CPU
# processname: /usr/sbin/powernowd

ARGS="-q"

. /lib/lsb/init-functions
