# compilers/j2sdk/profile.d/j2sdk.rc
#CLASSPATH is necessary for oh so many programs
export CLASSPATH=/usr/lib/j2sdk1.4.1/jre/lib/rt.jar
