# zbeta/slim/init.d/slim
#!/bin/bash

# chkconfig: 5 95 05 
# description: slim is a desktop manager
# processname: slim
# config: /etc/slim.conf

ARGS="-d"

. /lib/lsb/init-functions $1

