# utils/rsyslog/systemd.d/rsyslog.service
[Unit]
Description=Syslog daemon

[Service]
ExecStart=/usr/sbin/rsyslogd -c4 -n
Restart=always

[Install]
WantedBy=syslog.target
