# ftp/proftpd/init.d/ftp
#!/bin/bash
#
# ftp	ProFTP ftp server
#
# chkconfig: 2345 80 20
# description: ftp daemon
# processname: proftpd
# config:      /etc/ftpusers
# config:      /etc/proftpd.conf
# pidfile:     /var/run/proftpd.pid

extended-status () {
	ftpcount
	ftpwho
}

usage () {
	echo "Usage: $0 {start|stop|reload|probe|restart|status|extended-status}"
}

. /lib/lsb/init-functions

