# devel/ivman/init.d/ivman
#!/bin/bash
#
# Startup script for ivman
#
# description: bla
# chkconfig: 345 99 01
# processname: ivman
# pidfile: /var/run/ivman.pid

. /lib/lsb/init-functions $1
