# kernel/irqbalance/init.d/irqbalance
#!/bin/bash
#
# irqbalance 	irq load balancing for SMP systems
#
# chkconfig: 345 40 60
# description: irq load balancing daemon for SMP systems
# processname: irqbalance

. /lib/lsb/init-functions

