# mail/courier-authlib/init.d/authdaemon
#!/bin/sh
#
# authdaemon	    Courier-Authdaemon
#
# chkconfig: 345 75 25
# description: Courier-Authdaemon
# processname: authdaemond
# pidfile: /var/spool/authdaemon/pid

case $1 in
'start' )
  if [ -f /etc/authlib/authdaemonrc ] ; then
    echo -n "Starting authdameon: "
    authdaemond start &&
    echo -e $RESULT_OK ||
    echo -e $RESULT_FAIL
  fi
  ;;

'stop' )
    echo -n "Stopping authdameon: "
    authdaemond stop &&
    echo -e $RESULT_OK ||
    echo -e $RESULT_FAIL
  ;;
esac
