# xfce4/core/libxfce4util/profile.d/xfce4.rc
PATH=$PATH:/opt/lunar/xfce4/bin
PKG_CONFIG_PATH=/opt/lunar/xfce4/lib/pkgconfig:$PKG_CONFIG_PATH

export PATH PKG_CONFIG_PATH
