# zbeta/libvirt/init.d/libvirtd
#!/bin/bash

# chkconfig: 345 99 01
# description: libvirt is a library controling the installed virtual machines
# processname: /usr/sbin/libvirtd
# sigreload: HUP

ARGS="  -d"

. /lib/lsb/init-functions $1


