# net/nsd/init.d/nsd
#!/bin/bash
#
# nsd             This shell script takes care of starting and stopping nsd
#
# chkconfig: 345 45 55
# description: nsd is a Domain Name Server (DNS) \
# that is used to resolve host names to IP addresses.
# processname: nsd

. /lib/lsb/init-functions

