# audio/nas/init.d/nasd
#!/bin/bash
#
# Startup script for nasd
#
# chkconfig:   345 95 03
# description: A network transparent, client/server audio transport system
# processname: nasd
# pidfile: /var/run/nasd.pid

# config settings
[ -f /etc/config.d/nasd ] && . /etc/config.d/nasd

. /lib/lsb/init-functions $1
