# utils/preload/init.d/preload
#!/bin/bash
#
# Startup script for dbus
#
# chkconfig: 2345 05 95
# description: speed up startup times for aplications
# processname: preload
# pidfile: /var/run/preload.pid

. /lib/lsb/init-functions $1
