# java/groovy/profile.d/groovy.rc
export GROOVY_HOME=/usr/share/java/groovy
export PATH=${PATH}:${GROOVY_HOME}/bin
