# zbeta/ntop/init.d/ntop
#!/bin/bash
#
# ntop		The network traffic probe that shows the network usage
#
# chkconfig: 345 90 10
# description: ntop daemon
# processname: ntop
# pidfile: /var/ntop/ntop.pid
# config: /etc/config.d/ntop.conf

if [ -f /etc/config.d/ntop.conf ]; then
  . /etc/config.d/ntop.conf
else
  ARGS="--daemon --user nobody"
fi

start () {
  default_start
}

reload () {
  default_restart
}

. /lib/lsb/init-functions

