# utils/incron/init.d/incrond
#!/bin/bash
#
# incrond       This shell script enables the incrond daemon
#
# chkconfig:    - 97 03
#
# description:  This is a daemon which works like cron, but handles filesystem events
#		instead of time periods.
#               and can send notifications via mail, dbus or syslog.
# processname:  incrond
# pidfile:	/var/run/incrond.pid


start() {
	echo -n $"Starting incrond: "
	incrond && echo -e $RESULT_OK || echo -e $RESULT_FAIL
	echo
}

. /lib/lsb/init-functions
