# utils/apcupsd/systemd.d/apcupsd.service
[Unit]
Description=APC UPS Monitor

[Service]
ExecStart=/usr/sbin/apcupsd -b

[Install]
WantedBy=multi-user.target
