# zbeta/eet-cvs/profile.d/e17-cvs.rc
export PATH=${PATH}:/opt/lunar/enlightenment/bin
