# virtual/libvirt/systemd.d/virtlogd.socket
[Unit]
Description=Virtual machine log socket

[Socket]
ListenStream=/var/run/libvirt/virtlogd-sock

[Install]
WantedBy=multi-user.target
