# zbeta/nuttcp/xinetd.d/nuttcp
# default: off
# description: network performance testing server (developed from ttcp/nttcp)
service nuttcp
{
	port		= 5000
	socket_type	= stream
	wait		= no
	user		= nobody
	server		= /usr/bin/nuttcp
	server_args	= -S
	log_on_failure	+= USERID
}

