# kde3/base/kdenetwork3/init.d/lisa
#!/bin/bash
#
#
# chkconfig: 35 94 06
#
# description:  LISa is a small daemon which is intended to run on end user
#               systems. It provides something like a "network neighbourhood",
#               but only relying on the TCP/IP protocol stack, no smb or
#               whatever.
#
# processname: /opt/lunar/kde/3/bin/lisa
#
# config:      /etc/lisarc

#PATH=/sbin:/bin:/usr/sbin:/usr/bin:/opt/lunar/kde/3/bin

. /lib/lsb/init-functions
