# chat/bitlbee/systemd.d/bitlbee.service
[Unit]
Description=BitlBee IRC/IM gateway
After=syslog.target

[Service]
ExecStart=/usr/sbin/bitlbee -F -n

[Install]
WantedBy=multi-user.target
