# zbeta/scala-bin/profile.d/scala.rc
export SCALA_HOME="/usr/scala/latest"
export PATH="$PATH:$SCALA_HOME/bin"
export MANPATH="$MANPATH:$SCALA_HOME/man"
