# games/wesnoth/systemd.d/wesnothd.service
[Unit]
Description=Wesnoth Server Daemon
After=network.target

[Service]
ExecStart=/usr/games/wesnothd

[Install]
WantedBy=multi-user.target
