# net/howl/init.d/mDNSResponder
#!/bin/bash
#
# mDNSResponder:       Starts the mDNSResponder Daemon
#
# chkconfig: 345 97 03
# description: This is a daemon which runs on Howl clients to perform \
#              Zeroconf service discovery on a network.  mDNSResponder must be $#              that use Howl for service discovery. \
#              mDNSResponder should not be running otherwise.
# processname: mDNSResponder
# config: /etc/howl/mDNSResponder.conf
# pidfile: /var/run/mDNSResponder.pid

. /lib/lsb/init-functions $1
