# net/howl/init.d/nifd
#!/bin/bash
#
# nifd:       Starts the nifd Daemon
#
# chkconfig: 345 33 67
# description: This is a daemon which runs on Howl clients to monitor the \
#              state of a network interface.  nifd must be running on systems \
#              that use autoipd and mDNSResponder to automatically obtain a  \
#              Link-Local IPv4 address and do Zeroconf service discovery. \
#              nifd should not be running otherwise.
# processname: nifd
# pidfile: /var/run/nifd.pid
# config: 

. /lib/lsb/init-functions $1

