# utils/ddclient/init.d/ddclient
#!/bin/bash
#
# ddclient      This shell script takes care of starting and stopping
#          	ddclient.
#
# chkconfig: 2345 65 35
# description: ddclient provides support for updating dynamic DNS services.
# processname: ddclient
# pidfile: /var/run/ddclient.pid

. /lib/lsb/init-functions
