# kde3/kdebase3/init.d/kdm
#!/bin/bash

# chkconfig: 5 95 05
# description: kdm is KDE's version of the xdm graphical login manager
# processname: /opt/lunar/kde/3/bin/kdm
# sigreload: HUP

. /lib/lsb/init-functions $1
