# libs/telldus-core/systemd.d/telldusd.service
[Unit]
Description=Telldus TellStick daemon

[Service]
Type=forking
ExecStart=/usr/sbin/telldusd

[Install]
WantedBy=multi-user.target
